# docker-files/php7/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
user=root
autostart=true
command=/usr/local/sbin/php-fpm
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
